# deploy/wc-health-check.timer
[Unit]
Description=WeChat Health Check Timer — 每9小时健康巡检
Requires=wc-health-check.service

[Timer]
OnCalendar=*-*-* 03:57
OnCalendar=*-*-* 15:57
RandomizedDelaySec=180
Persistent=true

[Install]
WantedBy=timers.target
